# sdcard-images-utils/nxp/patches/ubuntu_overlay/step1/usr/lib/systemd/system/network-gadget.path
[Unit]
Description=Enable network when path is available

[Path]
PathExists=/run/usb_gadget_ready
Unit=network-gadget.service

[Install]
WantedBy=multi-user.target