# sdcard-images-utils/nxp/patches/ubuntu_overlay/step1/usr/lib/systemd/system/uvc-gadget.path
[Unit]
Description=Enable uvc when path is available

[Path]
PathExists=/sys/kernel/config/usb_gadget/g1/functions/uvc.0
Unit=uvc-gadget.service

[Install]
WantedBy=multi-user.target
